# docker/tomcat/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:tomcat]
command=service tomcat7 start && tail -f /var/lib/tomcat7/logs/catalina.out
